# roles/user/templates/system.service
[Unit]
Description = User Service
[Service]
User=roboshop
Environment=MONGO=true
Environment=REDIS_HOST=redis-dev.roboshop.internalganesh
Environment=MONGO_URL="mongodb-dev.roboshop.internalganesh:27017/users"
ExecStart=/bin/node /home/roboshop/user/server.js
SyslogIdentifier=user

[Install]
WantedBy=multi-user.target